# scripts/com.loveisland.sentiment.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.loveisland.sentiment</string>

    <key>ProgramArguments</key>
    <array>
        <string>/Users/benraeder/Documents/Love Island/loveisland-sentiment/.venv/bin/python</string>
        <string>/Users/benraeder/Documents/Love Island/loveisland-sentiment/scripts/run_pipeline.py</string>
    </array>

    <key>WorkingDirectory</key>
    <string>/Users/benraeder/Documents/Love Island/loveisland-sentiment</string>

    <!-- Run at fixed clock times, every 3 hours. launchd also runs a missed
         time once when the Mac wakes, so sleep gaps get a catch-up run. -->
    <key>StartCalendarInterval</key>
    <array>
        <dict><key>Hour</key><integer>0</integer><key>Minute</key><integer>0</integer></dict>
        <dict><key>Hour</key><integer>3</integer><key>Minute</key><integer>0</integer></dict>
        <dict><key>Hour</key><integer>6</integer><key>Minute</key><integer>0</integer></dict>
        <dict><key>Hour</key><integer>9</integer><key>Minute</key><integer>0</integer></dict>
        <dict><key>Hour</key><integer>12</integer><key>Minute</key><integer>0</integer></dict>
        <dict><key>Hour</key><integer>15</integer><key>Minute</key><integer>0</integer></dict>
        <dict><key>Hour</key><integer>18</integer><key>Minute</key><integer>0</integer></dict>
        <dict><key>Hour</key><integer>21</integer><key>Minute</key><integer>0</integer></dict>
    </array>

    <key>RunAtLoad</key>
    <false/>

    <key>StandardOutPath</key>
    <string>/Users/benraeder/Documents/Love Island/loveisland-sentiment/outputs/cron.log</string>
    <key>StandardErrorPath</key>
    <string>/Users/benraeder/Documents/Love Island/loveisland-sentiment/outputs/cron.log</string>
</dict>
</plist>
